# backup/scripts/etc/systemd/system/urbackup-mc.service
[Unit]
After=urbackup-mc.service
[Service]
ExecStart=/etc/urbackup/startup.sh
[Install]
WantedBy=default.target
